# cron/ea-refresh.cron
# CRON_TZ supported by supercronic
CRON_TZ=UTC
# Evening (late night)
0 0  * * * curl -fsS -X POST "http://intel-api:4000/open-loops/refresh?hours=6&runType=evening" -H "Authorization: Bearer ${B_API_KEY}" >/dev/null

# Morning (early)
0 6  * * * curl -fsS -X POST "http://intel-api:4000/open-loops/refresh?hours=6&runType=morning" -H "Authorization: Bearer ${B_API_KEY}" >/dev/null

# Morning (midday-ish)
0 12 * * * curl -fsS -X POST "http://intel-api:4000/open-loops/refresh?hours=6&runType=morning" -H "Authorization: Bearer ${B_API_KEY}" >/dev/null

# Evening (early night)
0 18 * * * curl -fsS -X POST "http://intel-api:4000/open-loops/refresh?hours=6&runType=evening" -H "Authorization: Bearer ${B_API_KEY}" >/dev/null