# cron/ea-refresh.cron
# CRON_TZ supported by supercronic
CRON_TZ=Europe/London
0 8 * * * curl -s -X POST "http://intel-api:4000/open-loops/refresh?hours=12&runType=morning" >/dev/null
0 20 * * * curl -s -X POST "http://intel-api:4000/open-loops/refresh?hours=12&runType=evening" >/dev/null
